# package/base-files/files/etc/init.d/uuplugin
#!/bin/sh /etc/rc.common


START=99
start() {
    /bin/sh /etc/uu//uuplugin_monitor.sh &
}
