# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
chmod +x /etc/lingti/lingti
/etc/lingti/lingti &
exit 0

